# scripts/autobahn.service
[Unit]
Description=Autobahn Server
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/autobahn --config /etc/autobahn/config.toml
Restart=always
RestartSec=2
StartLimitBurst=5
StartLimitIntervalSec=60
User=root
Group=root

[Install]
WantedBy=multi-user.target